# services/dbus-modem-e3372/run
#!/bin/sh
exec 2>&1
exec /data/VenusOS-E3372h/dbus-modem-e3372.py -s /dev/$tty
